# mdk/client/extraRFS/etc/inittab
# See https://git.busybox.net/busybox/tree/examples/inittab
::sysinit:/linuxrc

::askfirst:-/bin/sh
tty2::askfirst:-/bin/sh
tty3::askfirst:-/bin/sh

# /sbin/getty invocations for selected ttys
#tty4::respawn:/sbin/getty 38400 tty4
#tty5::respawn:/sbin/getty 38400 tty5

::restart:/bin/init

# Stuff to do before rebooting
::ctrlaltdel:/bin/reboot
::shutdown:/bin/umount -r
::shutdown:/bin/swapoff -a

